# server/opik-server.service
[Unit]
Description=OPIK Server - Embedding + FAISS + Chatbot
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/opik-server
Environment=S3_BUCKET=s3-opik-bucket
Environment=AWS_REGION=ap-northeast-2
Environment=BEDROCK_MODEL=apac.anthropic.claude-sonnet-4-20250514
ExecStart=/usr/bin/python3 -m uvicorn opik_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
